# build-scripts.lightdm/configs/blfs-systemd-units-20180105/blfs/units/rpcbind.service
[Unit]
Description=RPC Bind Service
Requires=rpcbind.socket

[Service]
Type=forking
ExecStart=/sbin/rpcbind -w

[Install]
Also=rpcbind.socket
