# build-scripts.lightdm/configs/blfs-systemd-units-20180105/blfs/units/slapd.service
[Unit]
Description=OpenLDAP server daemon
After=network.target

[Service]
Type=forking
EnvironmentFile=/etc/default/slapd
ExecStart=/usr/sbin/slapd -u ldap -g ldap $SLAPD_OPTS

[Install]
WantedBy=multi-user.target
