# build-scripts/06-lightdm/configs/blfs-systemd-units-20180105/blfs/units/exim.service
[Unit]
Description=Exim Mail Daemon
After=network.target

[Service]
PIDFile=/run/exim.pid
ExecStart=/usr/sbin/exim -bdf -q30m
ExecReload=/bin/kill -HUP $MAINPID
PrivateTmp=true

[Install]
WantedBy=multi-user.target
