# build-scripts/06-lightdm/configs/blfs-systemd-units-20180105/blfs/units/krb5-kadmind.service
[Unit]
Description=Kerberos 5 Administration Server

[Service]
ExecStart=/usr/sbin/kadmind -nofork

[Install]
WantedBy=multi-user.target
