# build-scripts/06-lightdm/configs/blfs-systemd-units-20180105/blfs/units/nfs-server.target
[Unit]
Description=NFS File Server Daemons
Wants=rpc-mountd.service
Wants=rpc-statd.service

[Install]
WantedBy=multi-user.target
