# build-scripts/lightdm/configs/blfs-systemd-units-20180105/blfs/units/krb5-kpropd.service
[Unit]
Description=Kerberos 5 Propagation Server

[Service]
ExecStart=/usr/sbin/kpropd -S

[Install]
WantedBy=multi-user.target
